# wifi-lights.service
Description=wifi-lights service

[Service]
ExecStart=/usr/bin/python3  /home/gwendal/wifi-lights/src/scripts/running/switch_status_watcher.py
Restart=always
User=gwendal
Group=gwendal

[Install]
WantedBy=multi-user.target